# postInstallScripts/endeavourOS/endavoursInstaller.service
[Unit]
Description=Post reboot installer

[Service]
Type=oneshot
ExecStart=~/test.sh

