# postInstallScripts/endeavourOS/endavoursInstaller.service
[Unit]
Description=Post reboot installer

[Service]
Type=oneshot
ExecStart=bash ~/endeabourOSXmonad.sh

[Install]
WantedBy=multi-user.target
